# etc/systemd/system/vboxdrv.service.d/dkms-compile.conf
[Service]
ExecStartPre=/usr/sbin/dkms autoinstall -m vboxhost
