# etc/systemd/system/vboxdrv.service
[Service]
Type=oneshot
ExecStart=/bin/sh /usr/lib/virtualbox/vboxdrv.sh start

[Install]
WantedBy=multi-user.target
